# openclaw-config/ai.openclaw.mission-control-frontend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.openclaw.mission-control-frontend</string>

    <key>Comment</key>
    <string>OpenClaw Mission Control frontend (Next.js, port 3000)</string>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>ProgramArguments</key>
    <array>
      <string>/opt/homebrew/bin/node</string>
      <string>/Users/jleechan/projects_reference/openclaw-mission-control/frontend/node_modules/.bin/next</string>
      <string>dev</string>
      <string>--port</string>
      <string>3000</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jleechan/projects_reference/openclaw-mission-control/frontend</string>

    <key>StandardOutPath</key>
    <string>/tmp/mc-frontend.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/mc-frontend.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key>
      <string>/Users/jleechan</string>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
      <key>NEXT_PUBLIC_API_URL</key>
      <string>http://localhost:9010</string>
    </dict>
  </dict>
</plist>
