# openclaw-config/ai.openclaw.task-poller.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.openclaw.task-poller</string>

    <key>Comment</key>
    <string>OpenClaw Mission Control task poller — polls inbox and dispatches to agents</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/jleechan/projects_reference/openclaw-mission-control/scripts/run_task_poller.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>/Users/jleechan/projects_reference/openclaw-mission-control</string>

    <key>StandardOutPath</key>
    <string>/tmp/mc-task-poller.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/mc-task-poller.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
      <key>MISSION_CONTROL_BASE_URL</key>
      <string>http://localhost:9010</string>
      <key>MISSION_CONTROL_BOARD_ID</key>
      <string>PLACEHOLDER_BOARD_ID</string>
      <key>MISSION_CONTROL_TOKEN</key>
      <string>PLACEHOLDER_MC_TOKEN</string>
    </dict>
  </dict>
</plist>
