# openclaw-config/ai.openclaw.schedule.daily-checkin-9am.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.openclaw.schedule.daily-checkin-9am</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>@HOME@/.openclaw/run-scheduled-job.sh</string>
      <string>522e23a7-c7c1-41f2-b117-a3af05661578</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>@HOME@/.openclaw/logs/scheduled-jobs/daily-checkin-9am.out.log</string>
    <key>StandardErrorPath</key>
    <string>@HOME@/.openclaw/logs/scheduled-jobs/daily-checkin-9am.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>@HOME@/.bun/bin:@HOME@/.local/bin:@HOME@/bin:@HOME@/Library/pnpm:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
      <key>HOME</key>
      <string>@HOME@</string>
    </dict>
  </dict>
</plist>
